# deploy/nextwaves-dashboard.service
[Unit]
Description=Nextwaves Bot Dashboard (Streamlit)
After=network-online.target nextwaves-bot.service
Wants=network-online.target

[Service]
Type=simple
User=REPLACE_USER
Group=REPLACE_USER
WorkingDirectory=REPLACE_PATH
EnvironmentFile=REPLACE_PATH/.env
ExecStart=REPLACE_PATH/venv/bin/streamlit run dashboard/app.py --server.address 127.0.0.1 --server.port 8501 --server.headless true
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
